# deploy/sensors2mqtt-hwmon.service
[Unit]
Description=sensors2mqtt hwmon collector (sw-bb-25g)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/sensors2mqtt/bin/python -m sensors2mqtt.collector.hwmon
Restart=always
RestartSec=10
Environment=MQTT_HOST=ha.welland.mithis.com
Environment=POLL_INTERVAL=30

[Install]
WantedBy=multi-user.target
