# debian/sensors2mqtt-local.prerm
#!/bin/bash
set -e

if [ "$1" = "remove" ] || [ "$1" = "deconfigure" ]; then
    systemctl stop sensors2mqtt-local || true
    systemctl disable sensors2mqtt-local || true
    systemctl daemon-reload
fi

#DEBHELPER#
